# contrib/systemd/skcomms-api.service
# SKComms S2S API: sovereign realm-aware comms federation surface.
#
# Serves the FastAPI app (skcomms.api:app) on the loopback socket 127.0.0.1:9384
# per SOP.md section 5. The internet-facing surface is a Tailscale Funnel :443
# path-route mounted onto this loopback port (see the bootstrap script and SOP
# section 10); the socket itself is NEVER bound to a public interface.
#
# Install as a user unit:
#   mkdir -p ~/.config/systemd/user
#   cp contrib/systemd/skcomms-api.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now skcomms-api.service
#
# Secrets are never inlined here. Per-host overrides (bind host/port, agent name,
# any secret PATHS the app reads) live in the EnvironmentFile below, which is
# optional (leading '-'): the unit starts on defaults if the file is absent.

[Unit]
Description=SKComms S2S API (skcomms.api:app on 127.0.0.1:9384)
Documentation=https://github.com/smilinTux/skcomms
After=network-online.target
Wants=network-online.target
# Restart watchdog: never let systemd give up on the S2S rail. Without this,
# the default start-rate limit (5 starts in 10s) can park a fast crash loop
# in a permanent 'failed' state; with it, Restart=always + RestartSec=5 keep
# retrying forever, so the daemon always comes back after a crash. The nonce
# replay cache is durable (SQLite under skcomms state/), so the restart does
# NOT open a replay window on the public inbox.
StartLimitIntervalSec=0

[Service]
Type=simple
# Env defaults; override per-host by editing the EnvironmentFile (paths, not values).
EnvironmentFile=-%h/.config/skcomms/skcomms.env
Environment=SKCOMMS_API_HOST=127.0.0.1
Environment=SKCOMMS_API_PORT=9384
# Identity gate (coord 7d5344f2): refuse to start with no CapAuth private key.
# A keyless start comes up with a green socket and dead crypto, and a
# regenerated key TOFU-conflicts on every remote peer. Restore first:
#   skcomms identity restore <archive>     (SOP.md section 11)
# Keyless dev standups may set SKCOMMS_ALLOW_NO_IDENTITY=1 in the EnvironmentFile.
ExecStartPre=/bin/sh -c 'if [ "$${SKCOMMS_ALLOW_NO_IDENTITY:-0}" = "1" ]; then exit 0; fi; exec %h/.skenv/bin/skcomms identity check --strict'
ExecStart=%h/.skenv/bin/uvicorn skcomms.api:app --host ${SKCOMMS_API_HOST} --port ${SKCOMMS_API_PORT}
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
